# roles/traefik-cert-exporter/templates/systemd/traefik-cert-exporter.service.j2
[Unit]
Description=Traefik certificate exporter
Wants=network-online.target
After=network-online.target
OnFailure=unit-status-mail@%n.service

[Service]
Type=oneshot
ExecStart={{ traefik_cert_exporter_root }}/bin/run
