# roles/traefik/templates/systemd/traefik.service.j2
[Unit]
Description=Traefik
BindsTo=docker.service
After=docker.service
OnFailure=unit-status-mail@%n.service

[Service]
# docker stop/rm just in case the container has not been removed (e.g. the
# system crashed).
ExecStartPre=-/usr/bin/docker container stop '%p'
ExecStartPre=-/usr/bin/docker container rm '%p'

# May fail when the network already exists.
ExecStartPre=-/usr/bin/docker network \
                              create \
                              --ipv6 \
                              --subnet fd00:dead:beef::/48 \
                              web

# --add-host=host.docker.internal:host-gateway is required to access services
# on the host network.
ExecStart=/usr/bin/docker container \
                          run \
                          --rm \
                          --name '%p' \
                          --network web \
                          --publish 80:80 \
                          --publish 443:443 \
                          --volume /var/run/docker.sock:/var/run/docker.sock:z \
                          --privileged \
                          --volume '/etc/localtime:/etc/localtime:ro' \
                          --volume '{{ root }}/conf/acme/:/etc/traefik/acme/' \
                          --volume '{{ root }}/conf/traefik.yaml:/etc/traefik/traefik.yaml:ro' \
                          --volume '{{ traefik_dynamic_file_configuration }}:/etc/traefik/conf.d/:ro' \
                          --env 'LEGO_DISABLE_CNAME_SUPPORT=true' \
                          --env 'INWX_USERNAME={{ inwx.acme.username | mandatory }}' \
                          --env 'INWX_PASSWORD={{ inwx.acme.password | mandatory }}' \
                          --env 'INWX_PROPAGATION_TIMEOUT=0' \
                          --add-host=host.docker.internal:host-gateway \
                          --label traefik.enable=true \
                          --label 'traefik.http.routers.dashboard.entrypoints=https' \
                          --label 'traefik.http.routers.dashboard.rule=Host(`traefik.{{ network.domain | mandatory }}`) && (PathPrefix(`/api`) || PathPrefix(`/dashboard`))' \
                          --label 'traefik.http.routers.dashboard.service=api@internal' \
                          --label 'traefik.http.routers.dashboard.middlewares=traefik-auth' \
                          --label 'traefik.http.routers.dashboard.tls.certresolver=letsencrypt-trs-dns' \
                          --label 'traefik.http.routers.dashboard.tls.domains[0].main={{ network.domain }}' \
                          --label 'traefik.http.routers.dashboard.tls.domains[0].sans=*.{{ network.domain }}' \
                          --label 'traefik.http.middlewares.traefik-auth.basicauth.realm=Traefik' \
                          --label 'traefik.http.middlewares.traefik-auth.basicauth.users={{ traefik_dashboard_users }}' \
                          --health-cmd 'traefik healthcheck' \
                          {{ traefik_image | mandatory | quote }}

ExecStop=/usr/bin/docker container stop '%p'

# Restart 2 seconds after docker run exited with an error status.
Restart=always
RestartSec=2s

[Install]
WantedBy=multi-user.target
