# roles/commafeed/templates/systemd/commafeed.service.j2
[Unit]
Description=Commafeed
BindsTo=docker.service
After=docker.service
OnFailure=unit-status-mail@%n.service

[Service]
WorkingDirectory={{ commafeed_root }}

ExecStartPre=/usr/bin/docker compose config
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose stop

ExecReload=/usr/bin/docker compose up --detach

[Install]
WantedBy=multi-user.target
