# roles/dns-server/templates/systemd/dns-server.target
[Unit]
Description=DNS Server
Wants=adguard.service dnsmasq.service dnscrypt-proxy.service

[Install]
WantedBy=multi-user.target
