# roles/immich/templates/systemd/data-immich-app.mount.j2
[Unit]
Description=Mount {{ immich_root }}/app

[Mount]
What=/fs/data/immich/app
Where={{ immich_root }}/app
Type=none
Options=bind
