# roles/immich/templates/systemd/data-immich-app.mount.j2
[Unit]
Description=Mount {{ root }}/app

[Mount]
What=/fs/data/immich/app
Where={{ root }}/app
Type=none
Options=bind
