# roles/docker/templates/systemd/buildx-prune/docker-buildx-prune.timer.j2
[Unit]
Description=Docker Buildx Prune Timer
BindsTo=docker.service

[Timer]
OnCalendar=4:30
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
