# roles/docker/templates/systemd/buildx-prune/docker-buildx-prune.service.j2
[Unit]
Description=Docker Buildx Prune
BindsTo=docker.service
After=docker.service
OnFailure=unit-status-mail@%n.service

[Service]
Type=oneshot

ExecStart=/usr/bin/docker buildx \
                          prune \
                          --all \
                          --force \
                          --filter until={{ docker_buildx_prune_keep_days * 24 | int }}h

[Install]
WantedBy=multi-user.target
