# roles/dns-server/templates/systemd/dnsmasq.service.d/override.conf
[Unit]
PartOf=dns-server.target

[Service]
WorkingDirectory={{ dns_server_root }}/dnsmasq

Type=simple
ExecStart=
ExecStart=/usr/sbin/dnsmasq --keep-in-foreground --conf-file=./dnsmasq.conf

ExecReload=/usr/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s
