# roles/docker/templates/systemd/image-prune/docker-image-prune.timer.j2
[Unit]
Description=Docker Image Prune Timer
BindsTo=docker.service

[Timer]
OnCalendar=4:00
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
