# roles/harbor/templates/systemd/data-harbor-app.mount.j2
[Unit]
Description=Mount {{ harbor_root }}/app

[Mount]
What=/fs/data/harbor/app
Where={{ harbor_root }}/app
Type=none
Options=bind
