# roles/gateway/templates/wireguard-resolve-endpoint/wireguard-resolve-endpoint.timer
[Unit]
Description=Resolve WireGuard peer endpoint IP address timer

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
