# roles/iobroker/templates/systemd/iobroker.service.j2
[Unit]
Description=ioBroker IoT platform
BindsTo=docker.service
After=docker.service adb-server.service
OnFailure=unit-status-mail@%n.service

[Service]
WorkingDirectory={{ iobroker_root }}

ExecStartPre=/usr/bin/docker compose config
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose stop

ExecReload=/usr/bin/docker compose up --detach

[Install]
WantedBy=multi-user.target
