# roles/iobroker/templates/systemd/iobroker-healthcheck.service.j2
[Unit]
Description=ioBroker IoT platform healthcheck
BindsTo=iobroker.service
OnFailure=unit-status-mail@%n.service

[Service]
ExecStart={{ root }}/healthcheck

[Install]
WantedBy=multi-user.target
