# roles/nextcloud/templates/systemd/data-nextcloud-backups.mount.j2
[Unit]
Description=Mount {{ nextcloud_root }}/backups

[Mount]
What=/fs/data/nextcloud/backups
Where={{ nextcloud_root }}/backups
Type=none
Options=bind
