# roles/borgmatic/templates/borgmatic.timer
[Timer]
OnCalendar=
OnCalendar=01:15
RandomizedDelaySec=10m
